# deploy/routerbot-tailscale.service
[Unit]
Description=Tailscale serve/funnel for RouterBot (ports 9420 and 10000)
After=tailscaled.service network-online.target
Wants=tailscaled.service
Before=routerbot.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=%ROUTERBOT_DIR%/scripts/tailscale-setup.sh
# Idempotent; safe on every boot

[Install]
WantedBy=multi-user.target
