# deploy/routerbot.service
[Unit]
Description=RouterBot OpenAI-compatible routing proxy
Documentation=file://%ROUTERBOT_DIR%/README.md
After=network-online.target tailscaled.service
Wants=network-online.target

[Service]
Type=simple
User=%ROUTERBOT_USER%
Group=%ROUTERBOT_GROUP%
WorkingDirectory=%ROUTERBOT_DIR%
EnvironmentFile=-/etc/routerbot/environment
ExecStart=%NODE_BIN% %ROUTERBOT_DIR%/src/server.js
Restart=always
RestartSec=5
# Avoid tight restart loops if node or config is broken
StartLimitIntervalSec=300
StartLimitBurst=10

# CLIs (claude, codex, cursor-agent, gemini) read auth from the service user's home
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
